# systemd/backup@.timer
[Unit]
Description=Timer for backup at %I

[Timer]
Persistent=yes
OnCalendar=daily
Unit=backup@%i.service

[Install]
WantedBy=multi-user.target
